# old/sentry-web.service
[Unit]
Description=Sentry Main Service
After=network.target
Requires=sentry-worker.service
Requires=sentry-cron.service

[Service]
Type=simple
User=sentry
Group=sentry
WorkingDirectory=/www/sentry
Environment=SENTRY_CONF=/etc/sentry
ExecStart=/www/sentry/bin/sentry run web
